# ansible/cfg/update-motd.d/95-ansible-command-exemple
#!/bin/sh

cat <<EOF

**************************************************************
* ANSIBLE     exemple
**************************************************************

EOF
echo "ansible all -i inventory.ini -m ping"
echo "ansible all -i inventory.ini -m shell -a \"uptime;free\""
